# rpi/config/web-server/bvsarui.service
[Unit]
Description=Start gunicorn

[Service]
ExecStart=/www/start.sh

[Install]
WantedBy=multi-user.target
